# Installer/payload/Library/LaunchAgents/io.fti.SplashBuddy.launch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>StartInterval</key>
    <integer>10</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>Label</key>
    <string>io.fti.SplashBuddy.launch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Application Support/SplashBuddy/SplashBuddy.launch.sh</string>
    </array>
</dict>
</plist>
